# debian/sixad.postinst
#!/bin/bash

set -e

#DEBHELPER#

update-rc.d sixad stop 10 0 1 6 .


ARCH=`uname -m`

case $ARCH in
    powerpc|ppc)
sed "s/Legacy=0/Legacy=1/" -i /etc/default/sixad
echo "Detected PowerPC[32] - Legacy driver will be used"
  ;;
esac
