# debian/sixad.postrm
#!/bin/bash

set -e

#DEBHELPER#

update-rc.d -f sixad remove
